# etc/rc.local
#!/usr/bin/env bash

# disable usb wakeup
# cat /proc/acpi/wakeup | grep XHCI
echo XHCI > /proc/acpi/wakeup
